# src/crontab
# run every 2 hours
* * * * *       /opt/solum-dashboard/launchpad-scripts/bugs.py
