# src/crontab
# run every 2 hours
0 */2 * * *       /opt/solum-dashboard/launchpad-scripts/bugs.py